# backend/config/crontab
0 10,22 * * * /app/spider.sh  >> /logs/spider_bid.log 2>&1